# deploy/lingyi-council.service
[Unit]
Description=灵依 议事厅守护进程
After=network.target

[Service]
Type=simple
User=ai
WorkingDirectory=/home/ai/LingYi
Environment=GLM_CODING_PLAN_KEY=a8525e2de71343de9de7f789a278df46.O6RtuHAKChQkVN0w
Environment=GLM_API_KEY=b10dcd0f050f4e16b2c795f4bcaa9d90.prLf6zofUgIvT3Kd
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 -m lingyi.cli council --interval 300
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
